# com.os3.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.os3.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>server</string>
    </array>

    <!-- 프로젝트 루트로 변경 필요 -->
    <key>WorkingDirectory</key>
    <string>/YOUR/PATH/TO/OS3</string>

    <!-- 부팅 시 자동 시작 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 크래시 시 자동 재시작 -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/os3-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/os3-server-error.log</string>
</dict>
</plist>
