# com.os2.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.os2.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>server</string>
    </array>

    <!-- Update this to your project root path -->
    <key>WorkingDirectory</key>
    <string>/YOUR/PATH/TO/project</string>

    <!-- Auto-start on boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Auto-restart on crash -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/os2-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/os2-server-error.log</string>
</dict>
</plist>
